# systemd/ml-monitor.service
[Unit]
Description=ML Drift Monitor and Auto Retrain Trigger
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/analytics-tool-v2
ExecStart=/usr/bin/env python3 /root/analytics-tool-v2/ml_monitoring.py --action monitor
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=120s

[Install]
WantedBy=multi-user.target
